# scripts/dubbing-api.service
[Unit]
Description=Dubbing Pipeline FastAPI Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/dubbing-mvp/python
Environment="PATH=/root/dubbing-mvp/python/venv/bin"
ExecStart=/root/dubbing-mvp/python/venv/bin/python /root/dubbing-mvp/python/app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

